# etc/cron.d/eswl-updater
# Update the whitelists once a day
20 2 * * *	root	/usr/local/bin/eswl-updater
